# config/supervisor-app.conf
# Uncomment if you don't use standalone flask app
#[program:nginx-app]
#command = /usr/sbin/nginx

[program:flask]
command = /usr/local/bin/python /code/app/main.py